# pkg/common/duetwebserver/usr/lib/systemd/system/duetwebserver.service
[Unit]
Description=Duet Web Server
Wants=duetcontrolserver.service
After=network.target duetcontrolserver.service

[Service]
WorkingDirectory=/opt/dsf/sd/www
ExecStart=/opt/dsf/bin/DuetWebServer
TimeoutStopSec=15
Type=simple

[Install]
WantedBy=multi-user.target
