# pkg/common/duetcontrolserver/usr/lib/systemd/system/duetcontrolserver.service
[Unit]
Description=Duet Control Server

[Service]
ExecStart=/opt/dsf/bin/DuetControlServer
TimeoutStopSec=15
Restart=always
RestartSec=5
Type=notify
User=dsf
Group=dsf
UMask=0002
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_SYS_TIME
AmbientCapabilities=CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_SYS_TIME

[Install]
WantedBy=sysinit.target
