# pkg/common/duetcontrolserver/usr/lib/systemd/system/duetcontrolserver.service
[Unit]
Description=Duet Control Server

[Service]
ExecStart=/opt/dsf/bin/DuetControlServer
TimeoutStopSec=15
Restart=always
RestartSec=5
Type=simple

[Install]
WantedBy=basic.target
